# templates/default/init/teamcity-agent.conf.erb
# teamcity-agent - TeamCity Agent
# THIS FILE IS MANAGED BY CHEF. ALL CHANGES WILL BE OVERWRITTEN!
#

description     "TeamCity Agent"

setuid <%= node["teamcity"]["user"] %>
setgid <%= node["teamcity"]["group"] %>

start on started teamcity-server
stop on stopped teamcity-server

pre-start exec /opt/teamcity/current/buildAgent/bin/agent.sh start
post-stop exec /opt/teamcity/current/buildAgent/bin/agent.sh stop
