# templates/default/init/teamcity-server.conf.erb
# teamcity-server - TeamCity Server
# THIS FILE IS MANAGED BY CHEF. ALL CHANGES WILL BE OVERWRITTEN!
#

description     "TeamCity Server"

env TEAMCITY_DATA_PATH=<%= node["teamcity-server"]["data_directory"] %>

start on filesystem or runlevel [2345]
stop on runlevel [!2345]

pre-start exec /opt/teamcity/current/TeamCity/bin/teamcity-server.sh start
post-stop exec /opt/teamcity/current/TeamCity/bin/teamcity-server.sh stop
